# examples/bls/extra/usr/lib/systemd/system/systemd-growfs-root.service.d/37-composefs.conf
# Make sure we grow the right root filesystem

[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-growfs /sysroot

